# node/supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D

[program:wildfly]
command=/opt/wildfly/bin/standalone.sh