# collector/supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D

[program:fluentd]
command=/etc/init.d/td-agent start
